# development/infra/containers/learninglocker/config/startup.conf
[program:apache2]
command=/home/apache.sh
numprocs=1
autostart=true
autorestart=true